# systemd/system/envirocam-upload.timer
[Unit]
Description=Envirocam Upload Service

[Timer]
#OnBootSec=120
OnBootSec=180

[Install]
WantedBy=timers.target
