# systemd/system/envirocam-photos.timer
[Unit]
Description=Envirocam Photos Service

[Timer]
OnBootSec=150

[Install]
WantedBy=timers.target
