# systemd/system/envirocam-detect-hang.timer
[Unit]
Description=Envirocam Hand Detection Service

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
