# systemd/system/envirocam-upload.timer
[Unit]
Description=Envirocam Upload Service

[Timer]
OnBootSec=120

[Install]
WantedBy=timers.target
