# systemd/system/envirocam-upload.timer
[Unit]
Description=Envirocam Upload Service

OnBootSec=60

[Install]
WantedBy=timers.target
