# systemd/system/envirocam-telemetry.service
[Unit]
Description=Envirocam Telemetry Service
After=envirocam-logging-socket.service

[Service]
Type=idle
User=pi

ExecStart=/usr/bin/bash /home/pi/dev/timelapse/scripts/saveTelemetry.sh
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
