# systemd/system/envirocam-upload.timer
[Unit]
Description=Envirocam Upload Service

[Timer]
OnBootSec=60

[Install]
WantedBy=timers.target
