# systemd/system/envirocam-photos.timer
[Unit]
Description=Envirocam Photos Service

[Timer]
OnBootSec=240

[Install]
WantedBy=timers.target
